# FluentFTPServer.Docker/glftpd/root/etc/service/xinetd/run
#!/bin/bash

# Start glfpd
if ! pid=$(pgrep xinetd);
then
	service xinetd start
fi
